# modules/osg-ce/files/gip_updater.cron
0 3 * * * cmssoft /opt/osg/app/etc/gip_updater.sh
